# 03_services/01_redis/configs/supervisor/redis.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
priority=10
user=redis
autostart=true
autorestart = true
stopsignal=QUIT
